# 7.3/crontab
* * * * * /usr/local/bin/php -r "file_put_contents(sys_get_temp_dir() . '/cron.log', date('Y-m-d H:i:s') . PHP_EOL, FILE_APPEND);"

# always left a blank line at the end of file.
